# deploy/weather-edge-hko-finalizer.timer
[Unit]
Description=Run WeatherEdge Hong Kong finalizer every day

[Timer]
OnCalendar=*-*-* 08:30:00 UTC
Persistent=true
Unit=weather-edge-hko-finalizer.service

[Install]
WantedBy=timers.target
